# Ansible/create_docker_start_control_hub_rmq_container_service/files/docker-start-control-hub-rmq-container.service
[Unit]
Description=Start the docker container specified.

[Service]
ExecStart=sudo docker start control-hub-rmq

[Install]
WantedBy=multi-user.target