# Ansible/Control Hub/create_docker_start_control_hub_rmq_container_service/files/docker-start-control-hub-rmq-container.service
[Unit]
Description=Start the docker container specified.

[Service]
Type=oneshot
ExecStart=sudo docker start control-hub-rmq
ExecStop=sudo docker stop control-hub-rmq
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target